# packages/zku/zku-1.0.23.tar.gz/zku-1.0.23/etc/systemd/system/zymkey-filesystem@.service
[Unit]
Description=Zymkey Filesystem Automounter %i
After=zkifc.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/zymkey filesystem_mount %i -a
ExecStop=/usr/local/bin/zymkey filesystem_unmount %i

[Install]
WantedBy=multi-user.target
