# packages/csr-azure-utils/csr_azure_utils-2.0.5.tar.gz/csr_azure_utils-2.0.5/auth-token.service
[Unit]
Description=Authentication Token service

[Service]
Type=simple
User=guestshell
Group=guestshell
WorkingDirectory=/home/guestshell
ExecStart=/bin/python /home/guestshell/.local/lib/python2.7/site-packages/csr_cloud/token_svr.py start
Restart=always
RestartSec=1
StandardOutput=null
StandardError=syslog

[Install]
WantedBy=multi-user.target
